# ansible/roles/ansible-role-haproxy/templates/etc/systemd/system/haproxy.service.d/limits.conf.j2
# {{ ansible_managed }}
[Service]
LimitNOFILE={{ haproxy_fd_limit }}
