# ansible/roles/ansible-role-certbot-hook/files/certbot.service
[Service]
ExecStart=
ExecStart=/usr/bin/certbot -q renew --deploy-hook /usr/local/bin/concat_certificates.sh
