# ansible/roles/ansible-role-offsite-backup/templates/etc/cron.d/restic_backup_veiled_haven
33 */6 * * * cstein RESTIC_PASSWORD={{offsite_backup_restic_password}} /tank/backup/veiled_haven/backup
